# debian/cron.d
*/1 * * * * root /usr/bin/watch-battery
