# debian/minos-tools.update-cd.cron.d
#This file is part of the minos-tools package
0 23 * * *  root /usr/bin/update-cd
