# debian/minos-desktop-tools.watch-battery.cron.d
#This file is part of the minos-desktop-tools package
*/1 * * * * root /usr/bin/watch-battery
