# debian/minos-tools.watch-battery.cron.d
#This file is part of the minos-tools package
*/1 * * * * root /usr/bin/watch-battery
