# debian/minos-tools.backup-mozilla.cron.d
#This file is part of the minos-tools-extra package
15 */4 * * * root /usr/bin/backup-mozilla
